# dockerizer/supervisor-laravel-tenants-schedule.conf
[program:laravel-tenants-schedule]
process_name=%(program_name)s_%(process_num)02d
command=/bin/bash -c "while true; do php artisan tenants:artisan schedule:run; sleep 60; done"
autostart=true
autorestart=true
startretries=3
numprocs=1
user=root
directory=/app/application
redirect_stderr=true
stdout_logfile=/var/log/laravel-tenants-schedule.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/laravel-tenants-schedule-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

